# services.d/nginx/run
#!/usr/bin/with-contenv sh
echo >&2 "starting nginx"
exec 2>&1
exec nginx -g "daemon off;"
