# services.d/tail/run
#!/usr/bin/with-contenv sh
echo >&2 "starting tail to follow log"
tail -f /var/log/nginx/access.log
